# debian/platform-modules-fishbone32.postinst
depmod -a
update-rc.d platform-modules-fishbone32 defaults
service platform-modules-fishbone32 start