# debian/platform-modules-questone2.postinst
systemctl enable sellog.timer
systemctl start sellog.timer
update-rc.d platform-modules-questone2 defaults
