# debian/platform-modules-fishbone48.postinst
depmod -a
update-rc.d platform-modules-fishbone48 defaults
service platform-modules-fishbone48 start