# debian/platform-modules-questone2.postinst
systemctl enable sellog.timer
systemctl start sellog.timer
update-rc.d platform-modules-questone2 defaults
service platform-modules-questone2 start
rm -f /etc/udev/rules.d/70-persistent-net.rules
