# debian/platform-modules-questone2.postinst
systemctl enable sellog.timer
systemctl start sellog.timer
service platform-modules-questone2 start
