# debian/platform-modules-questone2.postinst
# Enable high temp logger
systemctl enable sellog.timer
systemctl start sellog.timer
